# crons/crontab
0 1 * * * . /.env && /backup-postgres.sh